# celeryd.conf
; ==================================
;  celery worker supervisor example
; ==================================

[program:celery]
; Set full path to celery program if using virtualenv
command=celery worker -A celery_tasks --loglevel=INFO

directory=/home/christian/Programmierung/python/UpLoad2/libconcoct
user=celery_worker
numprocs=1
stdout_logfile=/home/christian/Programmierung/python/UpLoad2/libconcoct/log/worker.log
stderr_logfile=/home/christian/Programmierung/python/UpLoad2/libconcoct/log/worker.log
autostart=false
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

; When resorting to send SIGKILL to the program to terminate it
; send SIGKILL to its whole process group instead,
; taking care of its children as well.
killasgroup=true

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=998
